# systemd/.config/systemd/user/kdeindicator.service
[Unit]
Description=KDE Connect indicator
PartOf=graphical-session.target

[Service]
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/kdeconnect-indicator
Restart=always

[Install]
WantedBy=i3-session.target